# scripts/for_macine_parsonal/matsuichi/resources/etc/cron.d/ddns_sync
# DDNSの更新を15分ごとに実行する。
# m h dom mon dow user  command
*/15 * * * *    root    /usr/local/bin/ddns_sync.sh